# target/linux/ramips/base-files/etc/uci-defaults/14_rfid
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#
. /lib/functions.sh
. /lib/ramips.sh
. /lib/functions/system.sh


board=$(ramips_board_name)

case "$board" in
wrtnode | \
wt3020 | \
ur-326n4g | \
zbt-wa05)
	/etc/init.d/neard enable
	/etc/init.d/neard start
	;;
esac
