# build_dir/target-mipsel_24kec+dsp_uClibc-0.9.33.2/root-ramips/usr/lib/opkg/info/dbus.list
/usr/bin/dbus-launch.real
/usr/lib/dbus-1/dbus-daemon-launch-helper
/usr/sbin/dbus-daemon
/usr/bin/dbus-uuidgen
/etc/init.d/dbus
/usr/bin/dbus-launch
/etc/dbus-1/session.conf
/etc/dbus-1/system.conf
