# files/etc/init.d/disable
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org

START=99




start() {
        sleep 10

        lua /usr/bin/disable.lua

}

stop() {

        killall -9 relayr
}

