# files/etc/hotplug.d/iface/99-disable
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

sleep 10
lua /usr/bin/disable.lua

